# udev/rules.d/01-monitors.rules
ACTION=="change", SUBSYSTEM="drm", RUN+="sh -c 'echo hi > /tmp/mon'"
